# systemd/jetson/magicbox-web.service
[Unit]
Description=Magic Box Web Control Interface
Documentation=https://github.com/michihitoTakami/gpu_os

# Dependencies
After=network.target gpu-upsampler.service
Requires=gpu-upsampler.service

[Service]
Type=simple

# Run as non-root user
User=magicbox
Group=magicbox

# RTP はコードとして残すが、I2Sメイン運用ではデフォルトで無効化する
Environment=MAGICBOX_ENABLE_RTP=false
Environment=MAGICBOX_RTP_AUTOSTART=false

# Web UI bind defaults:
# - USB Gadget (usb0) only by default to avoid accidental exposure on Wi-Fi/Ethernet.
# - Override explicitly if needed (LAN内運用のみ。インターネット公開は禁止/非推奨):
#   MAGICBOX_WEB_HOST=0.0.0.0, MAGICBOX_WEB_PORT=80 など。
Environment=MAGICBOX_WEB_HOST=192.168.55.1
Environment=MAGICBOX_WEB_PORT=80

# Execution
WorkingDirectory=/opt/magicbox
ExecStart=/opt/magicbox/venv/bin/uvicorn web.main:app \
    --host ${MAGICBOX_WEB_HOST} \
    --port ${MAGICBOX_WEB_PORT} \
    --workers 1 \
    --log-level warning

# Restart policy
Restart=always
RestartSec=3

# Resource limits
MemoryMax=512M
CPUWeight=50

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/opt/magicbox/data /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=magicbox-web

[Install]
WantedBy=multi-user.target
