# systemd/jetson/totton-audio-gadget.service
[Unit]
Description=Totton Audio Project USB Composite Gadget (UAC2 + ECM)
Documentation=https://github.com/michihitoTakami/totton_audio/docs/jetson/usb-gadget/

# Start early in boot sequence, before network
DefaultDependencies=no
Before=network-pre.target
After=local-fs.target systemd-modules-load.service
Requires=local-fs.target

# Only start if USB Device Controller exists
ConditionPathExists=/sys/class/udc

[Service]
Type=oneshot
RemainAfterExit=yes

# Main commands
ExecStart=/opt/totton_audio/scripts/setup-usb-gadget.sh start
ExecStop=/opt/totton_audio/scripts/setup-usb-gadget.sh stop
ExecReload=/opt/totton_audio/scripts/setup-usb-gadget.sh restart

# Recovery on failure
Restart=on-failure
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
