# systemd/gpu-upsampler.service.in
[Unit]
Description=GPU Audio Upsampler (CUDA 16x to ALSA)
Documentation=https://github.com/michihitoTakami/michy_os

# ALSA/TCP-only path
After=network-online.target sound.target
Wants=network-online.target

[Service]
# Service type depends on libsystemd availability at build time:
# - Type=notify (HAVE_SYSTEMD=ON): Daemon sends READY=1/WATCHDOG=1
# - Type=simple (HAVE_SYSTEMD=OFF): No systemd notification
Type=@SYSTEMD_SERVICE_TYPE@
@SYSTEMD_NOTIFY_ACCESS@
WorkingDirectory=@CMAKE_INSTALL_FULL_DATADIR@/gpu_upsampler
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/gpu_upsampler_alsa
ExecReload=/bin/kill -HUP $MAINPID

@SYSTEMD_WATCHDOG_SEC@

# Graceful shutdown timeout: Allow 3 seconds for fade-out and cleanup
# After this timeout, systemd will send SIGKILL
TimeoutStopSec=3

# Restart policy
Restart=on-failure
RestartSec=3

# OOM protection: This is a critical audio service, protect from OOM killer
# -1000 = never kill, -900 = very unlikely to kill
OOMScoreAdjust=-900

# Security hardening (compatible with CUDA)
# Note: MemoryDenyWriteExecute=yes breaks CUDA, so omitted
LockPersonality=yes
NoNewPrivileges=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Resource management
Slice=session.slice

[Install]
# Standard system service (ALSA/TCP)
WantedBy=multi-user.target
