# systemd/raspberry_pi/raspi-control-api.service
[Unit]
Description=Magic Box Pi Control API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/magicbox/gpu_os
Environment=RPI_CONTROL_BIND_INTERFACE=usb0
Environment=RPI_CONTROL_PORT=8081
Environment=RPI_CONTROL_CONFIG_PATH=/opt/magicbox/gpu_os/raspberry_pi/usb_i2s_bridge/usb-i2s-bridge.env
Environment=RPI_CONTROL_STATUS_PATH=/var/run/usb-i2s-bridge/status.json
Environment=RPI_CONTROL_RESTART_CMD=systemctl restart usb-i2s-bridge.service
ExecStart=/usr/bin/python3 -m raspberry_pi.control_api
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
